# codeseek.service
[Unit]
Description=CodeSeek Application
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=codeseek
WorkingDirectory=/opt/codeseek/backend
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=codeseek
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target